# internal/systemd/units/lerd-autostart.service
[Unit]
Description=Lerd Autostart
After=default.target

[Service]
ExecStart=%h/.local/bin/lerd start
Restart=no

[Install]
WantedBy=default.target
